# contrib/systemd/watchdogd.service.in

[Unit]
Description=watchdog daemon
After=syslog.target

[Service]
Type=notify
NotifyAccess=all
ExecStart=@prefix@/sbin/watchdogd -F
NotifyAccess=all
OOMScoreAdjust=-1000
ExecReload=/bin/kill -1 $MAINPID
WatchdogSec=60
[Install]
WantedBy=multi-user.target
