# contrib/systemd/watchdogd.service.in

[Unit]
Description=watchdog daemon
After=syslog.target

[Service]
Type=simple
ExecStart=@prefix@/sbin/watchdogd -F
ControlGroup=cpu:/
NotifyAccess=main
OOMScoreAdjust=-1000

[Install]
WantedBy=multi-user.target
